# environment/etc/profile.d/vehicle-env.sh
# cam-server, udp-splitter, media-server
export VID_INPUT_PORT=5432
export VID_CLIENT_PORT=5400
export VID_CMD_PORT=7070
export MEDIA_SERVER_DIR=/home/pi/media
export MEDIA_SERVER_PORT=5595

